# cookbooks/swift/files/default/etc/systemd/system/grafana.service
[Unit]
Description=Grafana Server
After=multi-user.target

[Service]
ExecStart=/usr/local/bin/grafana

[Install]
WantedBy=multi-user.target
