# cookbooks/swift/files/default/etc/systemd/system/pykmip-server.service
[Unit]
Description=PyKMIP server
After=multi-user.target

[Service]
ExecStart=/usr/local/bin/pykmip-server

[Install]
WantedBy=multi-user.target
